# configs/blockchain-exporter.service
[Unit]
Description=F1R3FLY Blockchain Metrics Exporter
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu
ExecStart=/usr/bin/python3 /home/ubuntu/scripts/blockchain_metrics_exporter.py
Restart=always
RestartSec=10
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target